# wrtiot/config/wrtiot.init
#!/bin/sh /etc/rc.common
# Copyright (c) 2012 Edwin Chen <edwin@dragino.com>
# Copyright (c) 2014 Mikeqin <Fengling.Qin@gmail.com>

START=98

start()
{
	stty -F /dev/ttyUSB0 115200
	/usr/sbin/IoTd 
}

stop()
{
	/usr/sbin/IoTd -k
	rm -rf /var/run/IoTd
}
